# GenWatch-main/deploy/systemd/genwatch.service
[Unit]
Description=GenWatch — Generac H-100 monitoring & control
Documentation=https://github.com/sethmorrowsoftware/genwatch
After=network-online.target
Wants=network-online.target

[Service]
# Type=notify: the app sends sd_notify(READY=1) once startup is complete
# and pings sd_notify(WATCHDOG=1) every WatchdogSec/2 while healthy. If
# the poller hangs and watchdog pings stop, systemd SIGKILLs and restarts
# the unit — protecting against pymodbus / serial driver deadlocks.
Type=notify
NotifyAccess=main

User=genwatch
Group=dialout
WorkingDirectory=/opt/genwatch

# Config + secrets live in /etc/genwatch.
Environment=GENWATCH_CONFIG_PATH=/etc/genwatch/config.yaml
Environment=GENWATCH_DATA_DIR=/var/lib/genwatch
Environment=GENWATCH_UI_DIR=/usr/share/genwatch/ui
# Binds the control API (which includes generator + ATS command endpoints)
# on ALL interfaces. Auth (bcrypt+JWT), CSRF Origin checks, and a login rate
# limiter defend it, but it is plain HTTP — put a TLS reverse proxy (Caddy)
# or Tailscale in front for anything beyond a trusted wired OT LAN. Set to
# 127.0.0.1 if a local reverse proxy terminates TLS.
Environment=GENWATCH_HOST=0.0.0.0
Environment=GENWATCH_PORT=8000
Environment=GENWATCH_LOG_LEVEL=info
Environment=PYTHONUNBUFFERED=1

ExecStart=/opt/genwatch/venv/bin/python -m genwatch serve

# Reliability: restart on any failure. The service is built to stay up
# even when the Modbus link is down (it surfaces a "comms lost" state in
# the UI and keeps retrying the link in the background), so a restart here
# almost always means a hard process crash — paced at RestartSec=30 so we
# don't burn the SD card with journal churn on a permanent fault.
#
# StartLimitIntervalSec=0 DISABLES the start rate limiter on purpose: this
# is an unattended SAFETY monitor, and a permanently-`failed` unit (the
# monitor silently offline) is more dangerous than a slow restart loop, so
# we never give up restarting. A kernel/USB wedge that a service restart
# can't fix is covered instead by the hardware watchdog
# (system.conf.d/10-genwatch-hwwatchdog.conf), which resets the whole Pi.
Restart=always
RestartSec=30
StartLimitIntervalSec=0

# Watchdog: systemd will SIGKILL and restart the unit if the app stops
# pinging sd_notify(WATCHDOG=1). The app pings every 30 s while the
# poller loop is running.
WatchdogSec=60s
TimeoutStartSec=60s
TimeoutStopSec=30s

# ── Hardening: least privilege ────────────────────────────────────────
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=true
PrivateTmp=yes
PrivateDevices=no              # need access to /dev/ttyUSB*, /dev/ttyAMA*
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
ProtectClock=yes
ProtectHostname=yes
# ProtectProc=invisible + ProcSubset=pid: hide other processes' /proc
# entries from this service and restrict /proc to PID directories only.
# Free defense — the service has no business enumerating peers.
ProtectProc=invisible
ProcSubset=pid
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
LockPersonality=yes
SystemCallArchitectures=native
# Drop all capabilities. Service runs as User=genwatch (unprivileged)
# so it inherits none anyway, but pinning the bounding set defends
# against a future config edit that adds User=root for "debugging".
CapabilityBoundingSet=
AmbientCapabilities=
# Newly-created files (DB rotations, future exports) get group-read at
# most — never world-readable. Pairs with the 0640 perms install.sh
# stamps on config.yaml.
UMask=0027
# Syscall allowlist. @system-service is systemd's canonical baseline
# for a userspace network service — covers everything FastAPI /
# uvicorn / pymodbus / sqlite3 need. The negation list strips a few
# classes that no Python web service has any business calling
# (mounting, swapping, rebooting, raw I/O, CPU emulation, deprecated
# syscalls). EPERM rather than KILL on violation so a stray call
# surfaces as an exception in the logs rather than a hard crash —
# easier to diagnose if @system-service later proves too narrow.
SystemCallFilter=@system-service
SystemCallFilter=~@mount @swap @reboot @raw-io @cpu-emulation @obsolete
SystemCallErrorNumber=EPERM
# Bias the OOM killer away from us — losing the monitor during memory
# pressure on the Pi is much worse than losing whatever else is
# competing. Negative values are "kill me last"; -200 puts us well
# behind every default-priority process.
OOMScoreAdjust=-200

# Read-write paths
ReadWritePaths=/var/lib/genwatch /etc/genwatch /var/log/genwatch

# ── Device access — covers all common Pi serial paths ────────────────
# char-tty*: covers any USB-serial adapter regardless of /dev/ttyUSB index
# ttyAMA10/ttyAMA0: Pi 5 / Pi 4 on-board UART
# serial0/serial1: stable raspi aliases for the on-board UART
# genwatch-modbus: udev symlink created by deploy/udev/99-genwatch-modbus.rules
DeviceAllow=char-ttyUSB rw
DeviceAllow=char-ttyACM rw
DeviceAllow=/dev/ttyAMA0 rw
DeviceAllow=/dev/ttyAMA10 rw
DeviceAllow=/dev/serial0 rw
DeviceAllow=/dev/serial1 rw
DeviceAllow=/dev/genwatch-modbus rw

# Resource caps — Pi 4 has 4 GB, Pi 5 has 4-8 GB. MemoryHigh is a soft
# throttle (reclaim pressure) hit first; MemoryMax is the hard OOM-kill
# backstop. A slow leak or large WS fan-out then gets throttled and logged
# rather than abruptly OOM-killing the monitor — which, with Restart=always,
# would otherwise crash-loop.
MemoryHigh=384M
MemoryMax=768M
TasksMax=128
LimitNOFILE=4096

# Logging to journal
StandardOutput=journal
StandardError=journal
SyslogIdentifier=genwatch

[Install]
WantedBy=multi-user.target
